# deploy/crm.service
[Unit]
Description=Mrowki Coloring CRM (Node.js)
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/mrowki-crm
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production
Environment=COOKIE_SECURE=1
# Uncomment & set:
# Environment=SITE_ORIGIN=https://mrowki-coloring.pl
# Environment=TELEGRAM_BOT_TOKEN=xxx
# Or read from EnvironmentFile:
EnvironmentFile=/opt/mrowki-crm/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mrowki-crm

[Install]
WantedBy=multi-user.target
